# crontab
*/5 * * * * . /etc/environment; python /opt/crane.py >> /logfile